# crontab
0 4 * * * root docker run --rm  reboot_router:latest
