# crontab
0 4 * * * sudo docker run --rm  reboot_router:latest
