# crontab
4 * * * * sudo docker run --rm  re:latest
